# task_tech_4h.xml
<?xml version="1.0" encoding="UTF-16"?><Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task"><RegistrationInfo><Description>VideoTurn technology every 4 hours</Description><Author>VideoTurn</Author></RegistrationInfo><Triggers><CalendarTrigger><StartBoundary>2026-06-06T00:00:00</StartBoundary><Repetition><Interval>PT4H</Interval><Duration>P1D</Duration></Repetition><ScheduleByDay><DaysInterval>1</DaysInterval></ScheduleByDay></CalendarTrigger></Triggers><Principals><Principal id="Author"><UserId>DESKTOP-86KH7UH\Xiangyu</UserId><LogonType>S4U</LogonType><RunLevel>HighestAvailable</RunLevel></Principal></Principals><Settings><MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy><StartWhenAvailable>true</StartWhenAvailable><RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable><ExecutionTimeLimit>PT3H</ExecutionTimeLimit></Settings><Actions Context="Author"><Exec><Command>D:\coding_202605051504_XY_Propose_Minutes\videoturn_202606011405_XY\run_daily_auto_publish.bat</Command><Arguments>--topic-profile tech</Arguments><WorkingDirectory>D:\coding_202605051504_XY_Propose_Minutes\videoturn_202606011405_XY</WorkingDirectory></Exec></Actions></Task>
